# decrypt.service
[Unit]

Descrption=DecryptService

[Service]
ExecStart=/home/cleraude/Python-Crypto/decrypt.sh

[Install]
WantedBy=multi.user.target
